# deb/DEBIAN/postinst
#!/bin/sh
set -e
if [ -x "/etc/init.d/shadowsocks" ]; then
    update-rc.d shadowsocks defaults >/dev/null
fi
